# gaboot.service
[Unit]
Description=Gaboot Shop Backend
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin

[Install]
WantedBy=multi-user.target
